# configs/tidebox.ini
[supervisord]
nodaemon=true
logfile = /tmp/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -D

[program:jackd]
command=/usr/bin/jackd -d dummy
user=tidal
priority=10
startsecs=10
environment=HOME="/home/tidal",USER="tidal"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/tidal/jackd.log
stderr_logfile=/home/tidal/jackd.error_log

[program:dirt]
command=/repos/Dirt/dirt
user=tidal
priority=30
startsecs=10
environment=HOME="/repos/Dirt",USER="tidal"
directory=/repos/Dirt
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/tidal/dirt.log
stderr_logfile=/home/tidal/dirt.error_log

[program:ffserver]
command=/usr/bin/ffserver -f /home/tidal/ffserver.conf
user=tidal
priority=40
startsecs=10
autostart=true
autorestart=true
redirect_stderr=false
stdout_logfile=/home/tidal/ffserver.log
stderr_logfile=/home/tidal/ffserver.error_log

[program:ffmpeg]
command=/usr/bin/ffmpeg -f jack -i ffmpeg -ac 2 http://localhost:8090/feed1.ffm
user=tidal
priority=50
startsecs=10
autostart=true
autorestart=true
redirect_stderr=false
stdout_logfile=/home/tidal/ffmpeg.log
stderr_logfile=/home/tidal/ffmpeg.error_log

[program:jackd_connect1]
command=/usr/bin/jack_connect dirt:output_0 ffmpeg:input_1
user=tidal
priority=500
startsecs=5
autostart=true
autorestart=false

[program:jackd_connect2]
command=/usr/bin/jack_connect dirt:output_1 ffmpeg:input_2
user=tidal
priority=500
startsecs=5
autostart=true
autorestart=false